# targets/alexandrite-workstation/overlay/etc/skel/.config/autostart/ibus-daemon.desktop
[Desktop Entry]
Type=Application
Name=ibus daemon
Comment=
Exec=ibus-daemon -r -d -x
StartupNotify=false
