# targets/alexandrite-workstation/overlay/alexandriteos-system/root/etc/systemd/system/alexandriteos-osinfo-fixer.path
[Unit]
Description=os-release fixer

[Path]
PathChanged=/usr/lib/os-release

[Install]
WantedBy=default.target
