# targets/alexandrite-workstation/overlay/alexandriteos-system/root/etc/systemd/system/alexandriteos-osinfo-fixer.service
[Unit]
Description=os-release fixer

[Service]
Type=oneshot
ExecStart=/usr/bin/cp -f /usr/share/alexandrite/os-release /usr/lib/
