# targets/alexandrite-workstation/overlay/alexandriteos-system/root/etc/systemd/system/autoremove-gnome-software-packagekit-pulgin.path
[Unit]
Description=remove gnome-software PackageKit plugins

[Path]
PathChanged=/usr/lib64/gnome-software/plugins-17/

[Install]
WantedBy=default.target
