# targets/astrust/overlay/etc/skel/.config/autostart/sysupd-notify.desktop
[Desktop Entry]
Type=Application
Name=system-update-notify
Comment=system update notify
Exec=bash /usr/libexec/sysupd-notify
StartupNotify=false
