# targets/alexandrite-workstation/overlay/alexandrite-shell/root/etc/skel/.config/autostart/after-install.desktop
[Desktop Entry]
Type=Application
Name=afterInstall
Comment=After Install prosess
Exec=bash /usr/share/alexandrite/afterInstall.sh
StartupNotify=false
